# systemd/codex-slot-relay.service
[Unit]
Description=Codex Slot Relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/codex-slot-relay
ExecStart=/usr/local/bin/codex-slot-relay --runtime-root /var/lib/codex-slot-relay serve
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
